# watcher/docker-watcher-watcher.conf
description "docker-watcher watcher"

start on runlevel [2345]
stop on runlevel [016]

chdir /opt/docker-watcher/watcher
exec  /opt/docker-watcher/watcher/watcher.py