# web/docker-watcher-web.conf
description "docker-watcher web"

start on runlevel [2345]
stop on runlevel [016]

chdir /opt/docker-watcher/web
exec  /opt/docker-watcher/web/web.py